# lockAppNet/etc/crontabs/root
*/15 8-23 * * *  /system/bin/lockAppNet
