# lockAppNet/etc/crontabs/root
*    0-11/3  *    *   *    /system/bin/lockAppNet
*    14-23/3 *    *   *    /system/bin/lockAppNet
*/10 13-14   *    *   6,0  /system/bin/lockAppNet
